# mrblib/mitamae/plugin/resource_executor/templates/supervisord_program.conf.erb
[program:<%= @name %>]
command = <%= @command %>
<% @settings.each do |key, value| -%>
<%= key %> = <%= value %>
<% end -%>
